# installer/deploybot-monitor.service
[Unit]
Description=DeployBot Monitor

[Service]
ExecStart=/opt/deploybot/monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target