# installer/deploybot-scheduler.service
[Unit]
Description=DeployBot Scheduler

[Service]
ExecStart=/opt/deploybot/scheduler.sh
Restart=always

[Install]
WantedBy=multi-user.target